# systemd/garden-agent.service
[Unit]
Description=garden-agent FastAPI service
After=network.target
StartLimitIntervalSec=120
StartLimitBurst=5

[Service]
Type=simple
User=YOUR_VM_USER
WorkingDirectory=/home/YOUR_VM_USER/apps/garden-agent
EnvironmentFile=/home/YOUR_VM_USER/apps/garden-agent/secrets.env
ExecStart=/home/YOUR_VM_USER/apps/garden-agent/.venv/bin/uvicorn garden.main:app --host 127.0.0.1 --port 8001
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=garden-agent

[Install]
WantedBy=multi-user.target
