# systemd/garden-backup.timer
[Unit]
Description=garden-agent nightly SQLite backup (02:00 UTC)

[Timer]
OnCalendar=*-*-* 02:00:00 UTC
Persistent=true
Unit=garden-backup.service

[Install]
WantedBy=timers.target
